# timebox-daemon.service
[Unit]
Description=TimeBox Evo notification daemon

[Service]
# TIMEBOX_ADDRESS (required), TIMEBOX_PIN (if your box asks),
# TIMEBOX_LATLON (optional, weather on the box's clock) and
# TIMEBOX_CLOCK (optional, default clock pages) live here:
EnvironmentFile=%h/.config/timebox/env
WorkingDirectory=%h/.local/share/timebox
ExecStart=%h/.local/share/timebox/.venv/bin/python timebox_daemon.py
# connect_le gives up after 5 attempts; let systemd retry until the box is
# reachable (boot, box off, adapter down).
Restart=on-failure
RestartSec=15

# Sandbox: read-only system and home, except the runtime dir (FIFO,
# PipeWire sockets) and divoom_protocol's session state file.
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%t -%h/.divoom-protocol

[Install]
WantedBy=default.target
